# initsystem/etc/cron.daily/shadow
#!/bin/sh

[ -x /usr/bin/grpck ] && /usr/bin/grpck -r 2>&1 | /usr/bin/logger -t grpck
[ -x /usr/bin/pwck ] && /usr/bin/pwck -qr 2>&1 | /usr/bin/logger -t pwck
