# initsystem/initsystem/etc/cron.daily/logrotate
#!/bin/sh

[ -x /usr/bin/logrotate ] && /usr/bin/logrotate /etc/logrotate.conf 2>&1 | /usr/bin/logger -t logrotate
