# initsystem/initsystem/etc/rc.d/alsa-utils
#!/bin/ash

rc_type="nodaemon"

rc_start_cmd="/usr/bin/alsactl restore"
rc_stop_cmd="/usr/bin/alsactl store"

. /etc/rc.d/rc.subr

rc_run "start stop" $@

# vim: set filetype=sh:
