# initsystem/initsystem/etc/rc.d/ntp
#!/bin/ash

rc_start_cmd="/usr/bin/ntpd -g -u ntp:ntp"

. /etc/rc.d/rc.subr

rc_run "start stop restart status" $@

# vim: set filetype=sh:
