# services/rc.d.syslog
#!/bin/ash

# This file is part of initsystem
#
# Copyright (C) 2014-2015 Valère Monseur (valere dot monseur at ymail dot com)
#
# This program is free software; you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation; either version 2 of the License, or
# (at your option) any later version.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License along
# with this program; if not, write to the Free Software Foundation, Inc.,
# 51 Franklin Street, Fifth Floor, Boston, MA 02110-1301 USA.

PATH=/opt/busybox
cd /
umask 022

. /etc/rc.subr

start() {
  pgrep -n -x /opt/busybox/syslogd 1>/dev/null 2>/dev/null && \
    { echo "syslogd is already running"; return 1; }

  rm -f /run/syslogd.pid

  /opt/busybox/syslogd -K
}

stop() {
  pkill -TERM -x /opt/busybox/syslogd
  pkill -CONT -x /opt/busybox/syslogd

  pwait /opt/busybox/syslogd 20 || \
    { echo "syslogd is still running"; return 1; }

  rm -f /run/syslogd.pid
}

restart() {
  stop
  start
}

status() {
  echo -n "syslogd is running as pid: "
  pgrep -x /opt/busybox/syslogd | tr '\n' ' '
  echo ""
}

rc_run "start stop restart status" $@
