# initsystem/etc/cron.weekly/archlinux-keyring-sync
#!/bin/sh

[ -x /usr/bin/archlinux-keyring-wkd-sync ] && /usr/bin/archlinux-keyring-wkd-sync 2>&1 | /usr/bin/logger -t archlinux-keyring-sync
