# initsystem/etc/cron.weekly/fstrim
#!/bin/sh

[ -x /usr/bin/fstrim ] && /usr/bin/fstrim -Av 2>&1 | /usr/bin/logger -t fstrim
