# initsystem/etc/cron.weekly/paccache
#!/bin/sh

[ -x /usr/bin/paccache ] && /usr/bin/paccache -r 2>&1 | /usr/bin/logger -t paccache
