# initsystem/etc/rc.d/syslog-ng
#!/bin/ash

rc_start_cmd="/usr/bin/syslog-ng"

. /etc/rc.d/rc.subr

rc_run "start stop restart reload status" $@

# vim: set filetype=sh:
